# immich-machine-learning.container
[Container]
Pod=immich.pod
ContainerName=immich_machine_learning
Image=ghcr.io/immich-app/immich-machine-learning:release
EnvironmentFile=immich.env
Volume=model-cache:/cache

[Service]
Restart=always
