# immich-machine-learning.container
[Container]
Pod=immich.pod
ContainerName=immich_machine_learning
EnvironmentFile=/etc/containers/systemd/immich/immich.env
Image=ghcr.io/immich-app/immich-machine-learning:release
Volume=model-cache:/cache

[Service]
Restart=always
